# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:opencode]
command=/usr/local/bin/opencode serve --hostname 0.0.0.0 --port 4096
directory=/app
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopsignal=TERM
stopwaitsecs=10

[program:agent-webui]
command=node /opt/agent-webui/packages/ide-app/dist/node/server.js
directory=/opt/agent-webui
autostart=true
autorestart=true
startsecs=3
environment=AGENT_WEBUI_PORT="3000",AGENT_WEBUI_OPENCODE_BASE="http://127.0.0.1:4096",AGENT_WEBUI_BASE_PATH="/webui",AGENT_WEBUI_DIRECTORY="/app",AGENT_WEBUI_MOCK="true",NODE_ENV="production"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopsignal=TERM
stopwaitsecs=10
